# arylic_multiroom.service
[Unit]
Description=Arylic Multiroom setup
After=network.target

[Service]
ExecStart=-/usr/local/bin/multiroom.sh

[Install]
WantedBy=default.target
